# ward.container
[Unit]
Description=Ward rootless pi and tmux environment
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
ContainerName=ward
Image=ward.build
Pull=never
HostName=ward
Network=host
UserNS=keep-id:uid=1000,gid=1000
User=1000
Group=1000
WorkingDir=/home/agent/Projects
Volume=%h/Projects:/home/agent/Projects
Volume=/opt:/opt:ro
Volume=%h/.pi:/home/agent/.pi
Volume=%h/.local/share/pnpm:%h/.local/share/pnpm
Volume=%h/.tmux.conf:/home/agent/.tmux.conf:ro
Volume=%h/.zshrc:/home/agent/.zshrc:ro
Volume=%h/.oh-my-zsh:/home/agent/.oh-my-zsh:ro
PidsLimit=4096
NoNewPrivileges=true
DropCapability=all
LogDriver=journald
StopTimeout=30

[Service]
Restart=on-failure
RestartSec=5s
TimeoutStopSec=45s
MemoryHigh=12G
MemoryMax=16G
TasksMax=4096

[Install]
WantedBy=default.target
